# btscript.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/btscript.sh

[Install]
WantedBy=default.target